# etc/supervisord.d/fullerite.ini
[program:fullerite]
command=/opt/qnib/fullerite/bin/start.sh
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true
