# etc/supervisord.d/qcollect.ini
[program:qcollect]
command=/opt/qnib/qcollect/bin/start.sh
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true
